# packages/rootfs/reference/orbitdb-relay-pinner/rootfs/orbitdb-relay-pinner-bootstrap-refresh.timer
[Unit]
Description=Periodic Aleph bootstrap refresh for orbitdb-relay-pinner

[Timer]
OnBootSec=20m
OnUnitActiveSec=6h
RandomizedDelaySec=15m
Persistent=true
Unit=orbitdb-relay-pinner-bootstrap-refresh.service

[Install]
WantedBy=timers.target
