# packages/rootfs/reference/orbitdb-relay-pinner/rootfs/orbitdb-relay-pinner-bootstrap-refresh.service
[Unit]
Description=Refresh Aleph bootstrap registration for orbitdb-relay-pinner
Wants=network-online.target
After=network-online.target orbitdb-relay-pinner.service
ConditionPathExists=/etc/default/orbitdb-relay-pinner.ready

[Service]
Type=oneshot
Environment=ENV_FILE=/etc/default/orbitdb-relay-pinner
Environment=DESCRIBE_SCRIPT=/usr/local/sbin/orbitdb-relay-pinner-describe.py
ExecStart=/usr/bin/python3 /usr/local/sbin/orbitdb-relay-pinner-bootstrap-refresh.py
StandardOutput=journal
StandardError=journal
