# packages/rootfs/reference/orbitdb-relay-pinner/rootfs/orbitdb-relay-pinner-bootstrap.service
[Unit]
Description=Temporary orbitdb-relay-pinner setup endpoint
Wants=network-online.target
After=network-online.target
ConditionPathExists=!/etc/default/orbitdb-relay-pinner.ready

[Service]
Environment=ENV_FILE=/etc/default/orbitdb-relay-pinner
Environment=READY_FILE=/etc/default/orbitdb-relay-pinner.ready
Environment=SERVICE_NAME=orbitdb-relay-pinner.service
Environment=BOOTSTRAP_SERVICE=orbitdb-relay-pinner-bootstrap.service
Environment=METADATA_FILE=/run/orbitdb-relay-pinner-setup-metadata.json
Environment=METADATA_ERROR_FILE=/run/orbitdb-relay-pinner-setup-metadata.error
ExecStart=/usr/bin/python3 /usr/local/sbin/orbitdb-relay-pinner-setup-server.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
